# crontab
0,5,10,15,20,25,30,35,40,45,50,55,60 * * * * /ect/init.d/mysql start
0,5,10,15,20,25,30,35,40,45,50,55,60 * * * * python /usr/monitor/system_status_monitoring.py
